# linux/config/systemd/user/mbsync.service
[Unit]
Description=Sync Mail

[Service]
Type=oneshot
ExecStart=/bin/env mbsync -aV
